# var/lib/opkg/info/busybox.list
/etc/busybox.links
/bin/sh
/etc/init/
/etc/init/busybox-cron.conf
/etc/init/busybox-udhcpc.conf
/bin/busybox
/etc/init.d/
/etc/syslog-startup.conf.busybox
/etc/init.d/hwclock.sh
/etc/init.d/syslog.busybox
