# var/lib/opkg/info/upstart-sysvcompat.list
/etc/default/rcS
/etc/init.d/rc
/etc/init.d/rcS
