# etc/init/alexad-prep.conf
start on starting alexad 

task

env PIDDIR=/var/run/alexad

script
    mkdir -p $PIDDIR
    chown daemon:root $PIDDIR
end script
