# etc/init/alexaspeechplayer-prep.conf
start on starting alexaspeechplayer 

task

env PIDDIR=/var/run/alexaspeechplayer

script
    mkdir -p $PIDDIR
    chown daemon:daemon $PIDDIR
end script
