# etc/init/audioplayerng-prep.conf
start on starting audioplayerng

task

env PIDDIR=/var/run/audioplayerng

script
    mkdir -p $PIDDIR
    chown daemon:daemon $PIDDIR
end script
