# etc/init/audioplayer.conf
description "Doppler Audio Player Daemon"
author      "Amazon.com, Inc."

start on started pulseaudio-main
stop on stopping pulseaudio-main

respawn

exec /usr/local/bin/audioplayer
