# etc/init/authd.conf
description "Doppler Device/Cloud Authentication Daemon"
author      "Amazon.com, Inc."

start on starting apps-normal-start
stop on stopping apps-normal-start

respawn

exec /usr/local/bin/authd
