# etc/init/avahi-daemon.conf
description	"avahi-daemon mDNS/DNS-SD daemon"
author          "Amazon.com Inc. "

start on started dbus
stop on stopping dbus

expect daemon
respawn

env opts="-D"

exec avahi-daemon ${opts}

