# etc/init/cleanup-crash-dumps.conf
description "Crash dump clean up script launcher"
author      "Amazon.com, Inc"

start on started loguploadd

task

exec /etc/init.d/cleanup-crash-dumps.sh

