# etc/init/crond-generate-crontab.conf
description "Boot-time creation of crontabs"
author      "Amazon.com, Inc."

start on starting busybox-cron

task

script
    /etc/init.d/crond-generate-crontab start
end script
