# etc/init/davsd.conf
description "Doppler Device Artifact Vending Service Daemon"
author      "Amazon.com, Inc."

start on started audio-shm-main
stop on stopping audio-shm-main

respawn

exec /usr/local/bin/davsd
