# etc/init/devcap-daemon.conf
description "Kindle Device Capabilities Daemon"
author      "Amazon.com, Inc."

start on started dbus
stop on late-shutdown

respawn

exec /usr/local/bin/devcap-daemon
