# etc/init/doppler-mode-main.conf
description "Abstract job to indicate that device is a Doppler"
author      "Amazon.com, Inc."

start on doppler-profile-main
stop on stopping apps-early-start
