# etc/init/dynconfig.conf
description "Kindle Dynamic Configuration Data Service"
author      "Amazon.com, Inc."

start on starting apps-normal-start
stop on stopping apps-normal-start

respawn

exec /usr/local/bin/dynconfig
