# etc/init/earconplayer-prep.conf
start on starting earconplayer

task

env PIDDIR=/var/run/earconplayer

script
    mkdir -p $PIDDIR
    chown daemon:daemon $PIDDIR
end script
