# etc/init/gadgetsd-prep.conf
start on starting gadgetsd

task

env PIDDIR=/var/run/gadgetsd

script
    mkdir -p $PIDDIR
    chown daemon:root $PIDDIR
end script

