# etc/init/hallod.conf
description "Doppler Alexa SIP Dameon"
author      "Amazon.com, Inc."

start on starting doppler-mode-main
stop on stopping doppler-mode-main

respawn

setgid audio

exec /usr/local/bin/hallod
