# etc/init/i18nd.conf
description "Doppler Localization Daemon"
author      "Amazon.com, Inc."

start on starting apps-early-start
stop on stopping apps-early-start

respawn

exec /usr/local/bin/i18nd
