# etc/init/idled.conf
description "Doppler Idle Watcher Daemon"
author      "Amazon.com, Inc."

start on starting doppler-mode-main
stop on stopping doppler-mode-main

respawn

exec /usr/local/bin/idled
