# etc/init/kernelmetrics.conf
description "Doppler Kernel Metrics Daemon"
author      "Amazon.com, Inc."

start on started apps-delayed-start

respawn

exec /usr/local/bin/kernelmetrics
