# etc/init/log-wifi-firmware-crash.conf
description "Wi-Fi Firmware Crash Log Launcher"
author      "Amazon.com, Inc."

start on wifi-crashed

task

exec /usr/local/bin/wifi-firmware-crash
