# etc/init/loguploadd.conf
description "Doppler Log Upload Daemon"
author      "Amazon.com, Inc."

start on started apps-delayed-start
stop on stopping apps-delayed-start


respawn

exec /usr/local/bin/loguploadd
