# etc/init/phd.conf
description "Kindle Phone Home Daemon"
author      "Amazon.com, Inc."

start on starting apps-normal-start
stop on stopping apps-normal-start

expect fork
respawn

exec /usr/local/bin/phd
