# etc/init/reboot.conf
description "Support for the reboot logic"
author      "Amazon.com, Inc."

start on runlevel [6]

post-start script
    sync
    reboot -d -f
end script
