# etc/init/save-the-date.conf
description "Before shutdown, save the current date/time"
author      "Amazon.com, Inc."

start on starting shutdown

task

script
    /usr/local/bin/save-the-date start
end script
