# etc/init/silence-dmesg.conf
# silence-dmesg
#
# dmesg is noisy. Shut it up when we're interactive.

start on runlevel [2345]

task

# Silence everything except emergencies (e.g. panic)
exec /bin/dmesg -n 1
