# etc/init/statd.conf
description "Doppler System Status Monitor Daemon"
author      "Amazon.com, Inc."

start on started apps-delayed-start
stop on stopping doppler-mode-main

respawn

exec /usr/local/bin/statd
