# etc/init/tinyrot.conf
description "Rotate logs on shutdown/reboot"
author      "Amazon.com, Inc."

start on starting shutdown

task

env PATH="/sbin:/bin:/usr/sbin:/usr/bin"

exec tinyrot
