# etc/init/todo.conf
description "Kindle TODO service"
author      "Amazon.com, Inc."

start on started apps-delayed-start
stop on stopping apps-delayed-start

respawn

exec /usr/local/bin/todo
