# etc/init/updateSoftwareVersion.conf
description "Updates software version inside dynconfig"
author      "Amazon.com, Inc."

start on starting apps-early-start

task

exec /usr/local/bin/updateSoftwareVersion
