# etc/init/whad-prep.conf
start on starting whad

task

env PIDDIR=/var/run/whad

script
    mkdir -p $PIDDIR
    chown daemon:daemon $PIDDIR
end script