# etc/udev/rules.d/81-ath6kl.rules
ACTION!="change", GOTO="ath6kl_end"
SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", RUN+="/sbin/initctl emit --no-wait wifi-crashed"
LABEL="ath6kl_end"