# etc/udev/rules.d/udev-bluetooth-hcierror.rules
ACTION=="change", SUBSYSTEM=="bluetooth", ENV{HCI_ERROR}=="FRAME_REASSEMBLY", RUN+="/usr/local/bin/bt_hci_error"
