# etc/network/if-up.d/ntpdate
#!/bin/sh

#
# CMD will trigger setting the time via NTP in a reliable way.
# This script is therefore no longer needed.
#
